# package/mongo-mms-agent_1.0-6/etc/init/mongo-mms-agent.conf
# mongo-mms-agent - serves monitoring information on our mongodb instances to 10gen's
#            monitoring system.  See http://www.10gen.com/mongodb-monitoring-service.

#
#   Copyright 2012 PayPerks Inc.
#
#   Licensed under the Apache License, Version 2.0 (the "License");
#   you may not use this file except in compliance with the License.
#   You may obtain a copy of the License at
#
#       http://www.apache.org/licenses/LICENSE-2.0
#
#   Unless required by applicable law or agreed to in writing, software
#   distributed under the License is distributed on an "AS IS" BASIS,
#   WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
#   See the License for the specific language governing permissions and
#   limitations under the License.
#

# 
#  CHANGELOG
#
#    2012-10-01 - Initial version (alex@mbird.biz)
#





description     "Mongo MMS-Agent DB Monitor"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

kill timeout 5

chdir /var/run/mongo-mms-agent

pre-start script
    if grep -q '@API_KEY@' /opt/mongo-mms-agent/settings.py || grep -q '@SECRET_KEY@' /opt/mongo-mms-agent/settings.py
    then
        logger -is -t "$UPSTART_JOB" "ERROR: api or secret key not set in settings.py!"
        stop ; exit 1
    fi
end script

# Old version: run as root
#exec /usr/bin/python /opt/mongo-mms-agent/agent.py

# New version: run as unprivileged user
exec su -s /bin/sh -c 'exec "$0" "$@"' mongomms -- /usr/bin/python /opt/mongo-mms-agent/agent.py

